# deploy/small-cap-stack.service
[Unit]
Description=small-cap-stack (IBKR Phase-1 tracker)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/small-cap-stack
# Never build on-box (#278): a build here competes with the live tracker for the box's 2 vCPU / 4 GB
# and leaves BuildKit cache behind. --pull missing fetches the image only if it isn't local yet (a
# fresh box / first boot); deploy.yml pins IMAGE_TAG in .env, which compose reads, so a reboot
# brings back the SAME image that was deployed rather than drifting to :latest.
ExecStart=/usr/bin/docker compose up -d --no-build --pull missing
ExecStop=/usr/bin/docker compose down
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
